# roles/package_management/files/paccache.timer
[Unit]
Description=Clean-up old pacman files

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=True

[Install]
WantedBy=multi-user.target
